# ansible/home-automation/roles/netdata/files/update_netdata.cron
# /etc/cron.d/update_netdata
0 0 * * 0 root /usr/src/netdata.git/netdata-updater.sh
