# ansible/home-automation/roles/gateway/files/update_dnsmasq.cron
# /etc/cron.d/update_dnsmasq
* * * * * root /usr/local/sbin/update_dnsmasq.sh
