# ansible/home-automation/roles/webserver/files/renew_certs.cron
# /etc/cron.d/renew_cert
0 0 * * 0 root certbot renew && systemctl reload apache2
