# linux/etc/cron.d/update_dnsmasq
* * * * * root /usr/local/sbin/update_dnsmasq.sh
